# profiles/openbox/desktop/etc/skel/.config/autostart/plank.desktop
[Desktop Entry]
Name=Plank
Comment=Simple Clean Dock
Type=Application
Exec=plank
Icon=plank
NoDisplay=false
Categories=Utility;
