# profiles/openbox/netbook/etc/skel/.config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Terminal=false
StartupNotify=true
Exec=conky -p1 -c ~/.config/conky/active
