# layers/music/etc/skel/.config/autostart/idesk.desktop
[Desktop Entry]
Name=idesk
GenericName=Desktop icons
Categories=Utility;
Icon=idesk
Exec=idesk
Terminal=false
Type=Application
