# profiles/standard/x11/etc/profile.d/20-java-options.sh
# Fix for java apps not using anit-aliased fonts
export _JAVA_OPTIONS='-Dawt.useSystemAAFontSettings=lcd'
