# profiles/xfce/netbook/etc/skel/.config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Name=conky
Terminal=false
StartupNotify=true
Exec=bash -c "conky -p 1 -c ~/.config/conky/active"
