# profiles/xfce/netbook/etc/skel/.config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Terminal=false
StartupNotify=true
Exec=bash -c "sleep 1 && exec conky -c ~/.config/conky/active"
